# pingone-import.service
[Unit]
Description=PingOne Import Tool
After=network.target

[Service]
Type=forking
User=<your-username>
WorkingDirectory=/path/to/pingone-import
ExecStart=/usr/bin/bash -c 'npm start'
ExecStop=/usr/bin/bash -c 'npm run stop:background'
Restart=on-failure
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=pingone-import

[Install]
WantedBy=multi-user.target